# supervisor-dash-deploy.conf
[program:dash-deploy]
command=python /path/to/dash-deploy/dash-deploy.py
autostart=true
autorestart=true
